# debian/nagios-check-clamscan/etc/cron.daily/clamscan
#!/bin/sh
#run by cron
logrotate -f /etc/logrotate_clamscan.conf


while read dir; 
do 
    clamscan -r -i --stdout $dir >> /var/log/clamav/clamscan.log
done < /etc/clamav/clamscan.conf

